# etc/profile.d/oneapi.sh
#/bin/bash
source /opt/intel/oneapi/setvars.sh --config=/opt/intel/oneapi/env.txt > /dev/null
